# systemd/pixflow-debug-api.service
[Unit]
Description=PixFlow Debug Host API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=maxymou
WorkingDirectory=/home/maxymou/PixFlow
Environment=DEBUG_HOST_API_BIND=0.0.0.0
ExecStart=/usr/bin/node /home/maxymou/PixFlow/scripts/pixflow-debug-api.js
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
